# stfu-notify-daemon.service
[Unit]
Description=The stfu-notify daemon

[Service]
Environment="DISPLAY=:0"
Type=simple
Restart=always
ExecStart=/usr/bin/stfu-notify-daemon
ExecReload=/usr/bin/kill -HUP $MAINPID
KillMode=process
TimeoutSec=180

[Install]
WantedBy=graphical-session.target


